# crontab
0 3 */7 * * docker-compose --file ~/ruby/docker-compose.jobs.yml run backup
