# crontab
@weekly docker-compose --file ~/ruby/docker-compose.jobs.yml run backup
